# .config/systemd/user/jupyter.service
[Unit]
Description=Jupyter

[Service]
ExecStart=jupyter lab --NotebookApp.notebook_dir=/home/ferres

[Install]
WantedBy=default.target

